# deploy/name.vadim.fomo-morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>name.vadim.fomo-morning</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>/Users/vadim/work/fomo/scripts/daily.sh</string>
    </array>

    <!-- 08:00 Europe/Zurich (system-local time, which is what launchd uses) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>8</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <!-- If the Mac is asleep at 08:00, launchd fires the job when it wakes. -->
    <key>RunAtLoad</key><false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/vadim</string>
        <key>PATH</key>
        <string>/Users/vadim/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/vadim/work/fomo</string>

    <key>StandardOutPath</key>
    <string>/Users/vadim/work/fomo/.logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/vadim/work/fomo/.logs/launchd.err.log</string>
</dict>
</plist>
